# root/etc/services.d/minetest/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc \
	minetestserver \
	--gameid "$GAME_NAME" \
	--config /defaults/minetest.conf \
	${CLI_ARGS} \
	--port 30000
