# root/etc/services.d/minetest/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc minetestserver \
	--gameid civtest_game \
	--config /config/.minetest/main-config/minetest.conf \
	${CLI_ARGS} \
	--port 30000
